# root/etc/services.d/flexget/run
#!/usr/bin/with-contenv bash
cd /config || exit
exec \
    s6-setuidgid abc flexget -c /config/config.yml daemon start
